# scripts/typematrix/85-typematrix.rules
ACTION=="add", ATTR{idProduct}=="2030", RUN+="/usr/bin/set_typematrix_dvp"
